# price_driven_switch/config/supervisord.conf
[supervisord]
nodaemon=true

[program:fastapi]
command=uvicorn price_driven_switch.__main__:app --host 0.0.0.0 --port 8080
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:streamlit]
command=streamlit run /app/price_driven_switch/frontend/1_📊_Dashboard.py --server.port 8501 --server.address 0.0.0.0
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
